# systemd/ground-station-api.service
[Unit]
Description=Ground Station API - Public Satellite Telemetry Gateway
Documentation=https://github.com/yourteam/satellite-telemetry
After=network-online.target telemetry-parser.service anomaly-detector.service
Requires=telemetry-parser.service anomaly-detector.service

[Service]
Type=simple
User=groundstation
Group=groundstation
WorkingDirectory=/opt/devops-satellite-telemetry/service-a

# Readiness gate: poll Service B and C health endpoints before starting
# This ensures Service A doesn't accept requests until dependencies are ready
ExecStartPre=/bin/bash -c 'echo "[$(date)] Waiting for dependencies to be ready..."; for i in $(seq 1 30); do if curl -s http://telemetry-parser:3002/health > /dev/null 2>&1 && curl -s http://anomaly-detector:3003/health > /dev/null 2>&1; then echo "[$(date)] Dependencies ready, starting Service A"; exit 0; fi; echo "[$(date)] Attempt $i/30 - dependencies not ready, retrying in 2s..."; sleep 2; done; echo "[$(date)] ERROR: Dependencies failed to become ready after 60s"; exit 1'
ExecStart=/opt/devops-satellite-telemetry/service-a/venv/bin/python /opt/devops-satellite-telemetry/service-a/app.py

Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ground-station-api

# Security hardening
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
